# code/etc/systemd/system/code2.service
[Unit]
Description=Sample Path Monitor Service

[Service]
Type=oneshot
ExecStart=logger "moo moo"
