# code/etc/systemd/system/code.service
[Unit]
Description=code sample servicce

[Service]
Type=oneshot
ExecStart=logger "moo moo"
