# code/etc/systemd/system/code2.path
[Unit]
Description=Sample Path Monitor Service

[Path]
DirectoryNotEmpty=/home/setup/code/upload
Unit=code2.service

[Install]
WantedBy=multi-user.target
